# systemd/bt-audio.path
[Unit]
Description=Monitor Bluetooth audio sink changes

[Path]
PathChanged=/tmp/bt-sink-trigger

[Install]
WantedBy=default.target
